# jarvis-bot.service
[Unit]
Description=AKPRO Telegram Bot
After=network.target

[Service]
Type=simple
User=red
WorkingDirectory=/home/red/AKPRO-bot
ExecStart=/usr/bin/node src/bot.js
Restart=always
RestartSec=10
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target